# deploy/iconography-guide.service
[Unit]
Description=Iconography Guide
After=network.target

[Service]
Type=simple
User=iconography
WorkingDirectory=/opt/iconography-guide
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

# Load secrets from a file not checked into git
EnvironmentFile=/etc/iconography-guide.env

# Hardening (light touch — adjust if you tighten later)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
